# swmodd/files/swmodd.init
#!/bin/sh /etc/rc.common

START=94
STOP=10

USE_PROCD=1
PROG=/usr/sbin/swmodd

start_service() {
	procd_open_instance swmodd
	procd_set_param command ${PROG}
	procd_set_param respawn
	procd_close_instance
}

stop_service() {
	service_stop ${PROG}
}

reload_service() {
	stop
	start
}
