# ieee1905/files/etc/init.d/ieee1905
#!/bin/sh /etc/rc.common

START=99
STOP=10

USE_PROCD=1
PROG=/usr/sbin/ieee1905d

validate_ieee1905_section()
{
	uci_validate_section ieee1905 ieee1905 "${1}" \
		'debug:bool:false' \
		'enabled:bool:false'
}

configure_ieee1905()
{
	local enabled debug

	validate_ieee1905_section ${1} || {
		echo "Validation of ieee1905 section failed"
		exit 1;
	}

	if [ ${debug} -eq 1 ]; then
		# Forward stdout of the command to logd
		procd_set_param stdout 1
		# Same for stderr
		procd_set_param stderr 1
	fi

	if [ ${enabled} -ne 1 ]; then
		exit 0;
	fi
}

configure_network()
{
	ebtables -L FORWARD|grep -iqE "1:80:C2:(0)+:(0)+:13.*-j.*DROP"
	if [ "$?" -ne 0 ]; then
		echo "Applying drop rule to drop pkts forwared by kernel to 1905.1 multicast mac"
		ebtables -A FORWARD -d 01:80:c2:00:00:13 -j DROP
	fi
}

start_service() {
	procd_open_instance ieee1905
	procd_set_param command ${PROG}
	configure_ieee1905 "ieee1905"
	configure_network
	procd_set_param respawn
	procd_close_instance
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "network"
	procd_add_reload_trigger "wireless"
	procd_add_reload_trigger "netmode"
}
