# wifimngr/files/etc/uci-defaults/ruleng.wps
uci -q set ruleng.wps=rule
uci -q set ruleng.wps.recipe='/etc/ruleng/wps.json'
