# port-management/files/etc/init.d/port_management
#!/bin/sh /etc/rc.common

START=15
USE_PROCD=1

. /lib/functions.sh
include /lib/network

configure_ports(){
	local cfg=$1
	local ifname enabled speed duplex autoneg eee pause

	config_get ifname $cfg ifname
	[ -d /sys/class/net/$ifname ] || return

	config_get enabled $cfg enabled 1
	config_get speed $cfg speed 1000
	config_get duplex $cfg duplex "full"
	config_get autoneg $cfg autoneg 1
	config_get eee $cfg eee 0
	config_get pause $cfg pause 0

	set_port_settings $ifname $enabled $speed $duplex $autoneg $eee $pause
}

start_service() {
	[ -f /lib/network/port.sh ] || return

	config_load ports
	config_foreach configure_ports ethport
	config_foreach configure_ports sfpport
}

service_triggers() {
	procd_add_reload_trigger ports
}
