# voice-client/files/etc/uci-defaults/99-voice_client
#!/bin/sh

uci -q get voice_client.call_filter0 >/dev/null || {

uci -q batch <<-EOT
	add voice_client call_filter
	rename voice_client.@call_filter[-1]=call_filter0
	set voice_client.call_filter0.block_foreign=0
	set voice_client.call_filter0.block_special_rate=0
	set voice_client.call_filter0.block_outgoing=0
	set voice_client.call_filter0.block_incoming=0
	commit voice_client
EOT

}

uci -q get voice_client.RINGING_STATUS >/dev/null || {

uci -q batch <<-EOT
	add voice_client ringing_status
	rename voice_client.@ringing_status[-1]=RINGING_STATUS
	set voice_client.RINGING_STATUS.status=0
	set voice_client.RINGING_STATUS.enabled=0
	set voice_client.RINGING_STATUS.shouldring=1
	commit voice_client
EOT

}

uci -q batch <<-EOT
        delete firewall.sip
        set firewall.sip=include
        set firewall.sip.path=/etc/firewall.sip
        set firewall.sip.reload=1
        commit firewall
EOT

exit 0

