# layer2interface/broadcom/etc/hotplug.d/dsl/10-devices
[ "$STATUS" == "up" -o "$STATUS" == "down" ] && /etc/init.d/dsl reload

