# uboot/files/uboot
#!/bin/sh /etc/rc.common

START=50

start() {
	. /lib/upgrade/uboot-upgrade
	uboot_upgrade /boot/uboot.img
}
