# icwmp/files/etc/uci-defaults/85-cwmp-set-userid
# Copy defaults by the factory to the cwmp UCI user section.


# Get factory base MAC.
baseMac=$(db -q get hw.board.basemac)

# Erase colon and space characters.
baseMac=${baseMac//:/}
baseMac=${baseMac// /}

# Caseing and fixed length string.
mac=$(printf "%12.12X" $((0x$baseMac)))

# Get system serial number.
serial=$(db -q get hw.board.serial_number)

# Get userid values
acs_userid=$(uci -q get cwmp.acs.userid)
cpe_userid=$(uci -q get cwmp.cpe.userid)

# Only set if they are empty
if [ -z "$acs_userid" ]
then
    uci -q set cwmp.acs.userid="${mac:0:6}-${serial}"
fi

if [ -z "$cpe_userid" ]
then
    uci -q set cwmp.cpe.userid="${mac:0:6}-${serial}"
fi

# No need for commit here, it is done by uci_apply_defaults().
