# netmode/files/etc/uci-defaults/ruleng.internet
uci -q set ruleng.internet=rule
uci -q set ruleng.internet.recipe='/etc/ruleng/internet.json'
